# Robust_AI_API_REST_Docker/docker/supervisor-conf.d/gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn app:app -b 0.0.0.0:5000 --log-level=debug
directory=/srv/

