# Robust_AI_API_REST_Docker/docker/supervisor-conf.d/redisQworker.conf
[program:redis-q-worker]
command=/usr/bin/python3 redisQworker.py
directory=/srv/

